# templates/etc/init/dspam-milter.conf
description	"dspam-milter"

start on runlevel [2345]
stop on runlevel [!2345]

setuid postfix

exec /usr/local/bin/dspam-milter --config /etc/dspam/dspam-milter.conf

